# root/etc/services.d/minetest/run
#!/usr/bin/with-contenv bash

declare world
# predefine some variables
world_part="--world "
world_set=$world

# set WORLD_STRING variable to use user-defined world or default world.
if [ -z "$world" ]; then
WORLD_STRING=" "
else
WORLD_STRING=($world_part$world_set)
fi

exec \
	s6-setuidgid abc minetestserver "${WORLD_STRING[@]}" \
	--config /config/.minetest/main-config/minetest.conf --port 30000
